# board/raspberrypi/mount.conf
# mount- Mount sfilesystems on boot
#
# This helper mounts filesystems in the correct order as the devices
# and mountpoints become available.

description	"Mount filesystems on boot"

start on startup

emits filesystems

script
	mount -o remount,rw /
	mkdir -p /dev/pts
	mkdir -p /dev/shm
	mount -a
	initctl emit filesystems
end script
