# board/raspberrypi/ntpdate.conf
# ntpdate
#
##

description "Sets the date/time from remote server"

start on started network

exec ntpdate -u 0.us.pool.ntp.org