# board/raspberrypi/tty1.conf
# tty1 - getty
#
# This service maintains a getty on tty1 from the point the system is
# started until it is shut down again.

start on runlevel [2]

stop on runlevel [!2345]

respawn
exec /sbin/getty 38400 tty1
