# package/upstart/network.conf
# networking
##

description "starts the network"

start on runlevel [2345]
stop on runlevel [!2345]

pre-start script
	/sbin/ifup -a
end script

post-stop script
	/sbin/ifdown -a
end script